# zfstools/zfs-auto-snapshot.service.in
[Unit]
Description=zfs-auto-snapshot __name__

[Service]
Type=oneshot
ExecStart=/usr/bin/zfs-auto-snapshot __name__ __count__
