# zfstools/zfs-auto-snapshot.timer.in
[Unit]
Description=Run zfs-auto-snapshot __name__

[Timer]
OnCalendar=__time__
Persistent=true

[Install]
WantedBy=timers.target
